# supervisord.conf
[supervisord]
nodaemon=true
user=root

[inet_http_server]
file=/var/run/supervisor.sock
chmod=0700
port = 9001
username = sz123
password = sz123

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username = sz123
password = sz123

[program:node]
command=yarn start:prod

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile=/var/log/nginx/error.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/nginx/access.log
stdout_logfile_maxbytes=10MB